# script_iptables.service
[Unit]
Description=Startup

[Service]
ExecStart=/home/sio-master/script_iptables.sh

[Install]
WantedBy=default.target
